# system/etc/profile.d/freeswitch-path.sh
#!/bin/sh
# Add FreeSWITCH binaries to PATH
export PATH="/usr/local/freeswitch/bin:$PATH"
