# system/etc/systemd/system/bella-ring.service
[Unit]
Description=Bella periodic ring -> main menu
After=freeswitch.service
Wants=freeswitch.service

[Service]
Type=oneshot
User=freeswitch
Group=freeswitch
ExecStart=/usr/local/freeswitch/scripts/bella-ring
